# canned_files/ld.so.preload
#CHROOT. we don't need anything to be preloaded for CHROOT-ed system
